# root/etc/init.d/hostname
#!/sbin/openrc-run
start() {
        ebegin "Set hostname"
        hostname $HOSTNAME
        eend $?
}
