# root/etc/init.d/i3wm
#!/sbin/openrc-run
depend() {
       after Xnest
}

run() {
   /usr/bin/i3 &
}

start() {
        ebegin "Starting i3wm window manager"
        run
        eend $?
}
