# root/etc/init.d/Xnest
#!/sbin/openrc-run
depend() {
       after system.init
}

run() {
        export DISPLAY=:0
        su -c "/usr/bin/Xnest :10" user &
}

start() {
        ebegin "Starting X server"
        run
        eend $?
}
