# root/etc/init.d/Xnest
#!/sbin/openrc-run
depend() {
       after system.init
}

run() {
        export DISPLAY=:0
        PARAM="-sync -name ShapeOS"
        if [ ! -z $GEOMETRY ]
        then
         PARAM="-geometry $GEOMETRY $PARAM"
        fi
        su -p  -c "/usr/bin/Xephyr -fullscreen $PARAM :10" user  > /var/log/X.log 2&>1 &
}

start() {
        ebegin "Starting X server"
        run
        eend $?
}
